# roles/setup_eos/files/etc/profile.d/eos.csh
U2FsdGVkX1/d+MviWsqyYkvn7HiHav4NX5yj5IkYRdCUsfpmVh7Y+Z/WmhRvF9xj
lAdnzn5xun8HMPFVUyLgCyvjnAKJkqVywCzC/LZByxMPko9V4Wh7p7cFO1kk8sJy
v4nbLWyom4rbdVJrgrfi6V1+IOF0caOUHj+iKNLazSVaHFEyVbS0zbtSxp+NapBC
