# roles/setup_eos/files/etc/profile.d/eos.sh
U2FsdGVkX1/yeXUluDguT5gtQhdLo2dJqimfpXZ3etigJZU0VGNijFILtMvIao4F
KYbSI3OJzoKzENDyk1j7uq4z7yll8lqR44p/zfC1J7l04XFeKdqKddItazuYcI8r
9HOLyrEhqQipuVoqGQrVz2I367xI8LYfRy2FoT/wGlDYJTjJEv8z5EkzD0KjmV4n
EDI5yf0b3bSdYqBGy9CmIs8s+7QV25qFazqR1UTY9WY=
